# systemd/user/mopidy.service
[Unit]
Description=Mopidy Daemon

[Service]
ExecStart=/usr/bin/mopidy
Type=simple

[Install]
WantedBy=default.target
